# build/install/docker/people.server/supervisord.conf
[supervisord]
nodaemon=true

[program:people]
directory=/var/www/products/ASC.People/server/
command=dotnet ASC.People.dll --urls=http://0.0.0.0:5004 --pathToConf=/app/onlyoffice/config/  --$STORAGE_ROOT=/app/onlyoffice/data/ --ENVIRONMENT=test
autostart=true
autorestart=true
